# docs/systemd/jarvis-tutor.service
[Unit]
Description=Aurelia Jarvis Tutor - Escuta Tudo 24/7
After=network.target audio.target
Wants=network.target

[Service]
Type=simple
User=will
WorkingDirectory=/home/will/aurelia
ExecStart=/home/will/aurelia/aurelia tutor --user=USER_ID --chat=CHAT_ID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jarvis-tutor
Environment="HOME=/home/will"
Environment="AURELIA_HOME=/home/will/aurelia"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /home/will/aurelia

# Logging
LogLevelMax=info

[Install]
WantedBy=multi-user.target
