# scripts/aurelia.service
[Unit]
Description=Aurelia Elite Daemon
After=network.target

[Service]
Type=simple
User=will
WorkingDirectory=/home/will/aurelia
ExecStart=/home/will/aurelia/aurelia-elite
Restart=always
RestartSec=5
StandardOutput=append:/home/will/.aurelia/logs/daemon.log
StandardError=append:/home/will/.aurelia/logs/daemon_error.log

[Install]
WantedBy=multi-user.target
